# daemon/alert/crontab/backup_affiliates.sh
#!/bin/bash

cd /var/www/html/daemon/alert
php ./backup_affiliates.php
