# daemon/alert/crontab/cron_get_konnektive_prospect.sh
#!/bin/bash

cd /var/www/html/APILotusCommercials/daemon/alert
php ./alert_konnektive_prospect.php


