# daemon/alert/crontab/cron_affiliate_cache_yesterday.sh
#!/bin/bash

cd /var/www/html/daemon/alert
php ./dashboard_affiliate_cache_yesterday.php


