# daemon/alert/crontab/cron_db_backup.sh
#!/bin/bash
cd /var/www/html/APILotusCommercials/daemon/alert
php ./db_backup.php
