# daemon/alert/crontab/cron_crm_password.sh
#!/bin/bash
cd /var/www/html/APILotusCommercials/daemon/alert
php ./crm_password.php
