# daemon/alert/crontab/cron_alert_status_schedule.sh
#!/bin/bash
cd /var/www/html/daemon/alert
php ./schedule.php
